# crontab/autoApt.job
#!/usr/bin/bash

apt-get update && apt-get upgrade -y && apt-get autoclean && apt-get autoremove -y