# crontab/autoApt.job
#!/usr/bin/env bash

## crontab -e
## 5 1 * * 1 /data/deploy/crontab/autoApt.job >> /dev/null 2>&1

apt update \
&& apt upgrade -y \
&& apt dist-upgrade -y \
&& apt autoremove -y \
&& apt autoclean